# dot_config/systemd/user/tldr-update.service
[Unit]
Description=Update tldr CLI cache
Documentation=https://tldr.sh/

[Service]
ExecStart=tldr --update
Type=oneshot
